# env-packagemanager/scheduled/package-database-update.service
[Unit]
Description=Update package cache

[Service]
ExecStart=/usr/bin/yaourt -Sy --aur
